# init/upstart/click-system-hooks.conf.in
description "Run Click system-level hooks"
author "Colin Watson <cjwatson@ubuntu.com>"

start on filesystem

task

exec @bindir@/click hook run-system
