# upstart/click-hook.conf
description "Run Click user-level hooks"
author "Colin Watson <cjwatson@ubuntu.com>"

start on starting xsession-init

exec click hook install-user
